# templates/influxdb3.container.j2
{# Managed by Ansible — metrics-aggregation role. Do not edit manually. #}
{# Instance: {{ influxdb3_instance_name }} #}
[Unit]
Description=InfluxDB 3 Core — instance {{ influxdb3_instance_name }}
After=network-online.target influxdb3-{{ influxdb3_instance_name }}-volume.service
Requires=influxdb3-{{ influxdb3_instance_name }}-volume.service

[Container]
Image={{ influxdb3_image }}
ContainerName=influxdb3-{{ influxdb3_instance_name }}
Volume=influxdb3-{{ influxdb3_instance_name }}.volume:/var/lib/influxdb3:Z
Network={{ influxdb3_podman_network }}
EnvironmentFile={{ influxdb3_env_file_path }}
# Disable auth on health/ping so Traefik and monitoring tools can probe without a token
Exec=influxdb3 serve --disable-authz health,ping

# Traefik discovery labels
Label=traefik.enable=true
Label=traefik.http.routers.influxdb3-{{ influxdb3_instance_name }}.rule=Host(`{{ influxdb3_instance_name }}.influxdb3.{{ influxdb3_traefik_domain }}`)
Label=traefik.http.routers.influxdb3-{{ influxdb3_instance_name }}.entrypoints={{ influxdb3_traefik_entrypoint }}
{% if influxdb3_traefik_entrypoint != 'web' %}
Label=traefik.http.routers.influxdb3-{{ influxdb3_instance_name }}.tls=true
{% endif %}
Label=traefik.http.services.influxdb3-{{ influxdb3_instance_name }}.loadbalancer.server.port={{ influxdb3_http_port }}

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec={{ influxdb3_startup_timeout_seconds }}

[Install]
WantedBy=default.target
